# opencarwings/rootfs/etc/services.d/opencarwings/run
#!/bin/sh
cd /opt/opencarwings || exit 1

echo "OpenCarwings add-on service starting..."

export PSQL_DATABASE=carwings
export PSQL_USER=carwings_user
export PSQL_PASSWORD=secure_password
export PSQL_DATABASE_HOST=localhost
export PSQL_DATABASE_PORT=5432
export REDIS_HOST=localhost
export REDIS_PORT=6379
export PYTHONPATH=/opt/opencarwings
export DJANGO_SETTINGS_MODULE=carwings.settings
export DEBUG=true
export SECRET_KEY=django-insecure-default-key-change-in-production

echo "Starting Django development server on port 8124..."

# Start Daphne ASGI server (handles both HTTP and TCP)
exec daphne -b 0.0.0.0 -p 8124 \
  --endpoint tcp:port=55230:interface=0.0.0.0 \
  carwings.asgi:application